# src/etc/supervisor.d/start_fail2ban.conf
[program:fail2ban]
command = /usr/bin/fail2ban-server -xf start
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
autorestart = true
